# systemd/fgip-tier1.timer
[Unit]
Description=Run FGIP Tier 1 agents every 6 hours

[Timer]
OnCalendar=*-*-* 00,06,12,18:00:00
Persistent=true
RandomizedDelaySec=600

[Install]
WantedBy=timers.target
