# systemd/fgip-tier3.timer
[Unit]
Description=Run FGIP Tier 3 meta-analysis agents daily (after tier 0)

[Timer]
# Run daily at 5:00 AM (after tier 0 at 2 AM)
OnCalendar=*-*-* 05:00:00
Persistent=true

[Install]
WantedBy=timers.target
